# overlay/etc/systemd/system/wifi-regdomain.service
[Unit]
Description=Unblock wifi rfkill
# cfg80211.ieee80211_regdom=US is already on the kernel cmdline. brcmfmac is
# self-managed and ignores cfg80211 anyway, so all we need here is to clear
# the soft-block that persists from the pristine image state.
After=systemd-rfkill.service
Before=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/rfkill unblock wifi

[Install]
WantedBy=multi-user.target
