# overlay/etc/systemd/system/xrce-dds-agent.service
[Unit]
Description=Micro XRCE-DDS Agent
After=network.target

[Service]
Type=simple
ExecStartPre=/usr/bin/python3 /etc/xrce-dds/serial-guard.py pre-start /dev/ttyAMA0 35
ExecStart=/opt/xrce-dds/bin/MicroXRCEAgent serial --dev /dev/ttyAMA0 -b 921600
ExecStopPost=/usr/bin/python3 /etc/xrce-dds/serial-guard.py post-stop
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
